# rkn-status/DEBIAN/postinst
#!/bin/bash

systemctl enable uwsgi-emperor
systemctl restart uwsgi-emperor
